# graphic/.config/autostart/conky.desktop
[Desktop Entry]
Exec=conky -d -c  ~/.config/conky/conkyrc
Name=conky
OnlyShowIn=LXQt;
Type=Application
Version=1.0
X-LXQt-X11-Only=true
