# graphic/.config/autostart/ibus.desktop
[Desktop Entry]
Exec=ibus-daemon -drx
Name=ibus
Type=Application
Version=1.0
X-LXQt-X11-Only=true
