# graphic/.config/autostart/udiskie.desktop
[Desktop Entry]
Exec=udiskie -A -t
Name=udiskie
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
X-LXQt-X11-Only=true
