# graphic/.config/autostart/volumeicon.desktop
[Desktop Entry]
Exec=volumeicon
Name=volumeicon
OnlyShowIn=LXQt;
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
X-LXQt-X11-Only=true
