# graphic/.config/autostart/xresources.desktop
[Desktop Entry]
Exec=xrdb ~/.Xresources
Name=xresources
Type=Application
Version=1.0
X-LXQt-X11-Only=true
fi=
